# Server/test.yorubanamedictionary.web.net.service
[Unit]
Description=Staging - Yoruba Name Dictionary Website .NET
After=network.target

[Service]
WorkingDirectory=/home/ec2-user/app/staging/yorubanamedictionary.web.net
ExecStart=/usr/bin/dotnet /home/ec2-user/app/staging/yorubanamedictionary.web.net/Website.dll
Restart=always
# Restart service after 30 seconds if the dotnet service crashes:
RestartSec=30
KillSignal=SIGINT
SyslogIdentifier=staging-yorubanamedictionary-website-net
User=ec2-user
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5001
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

[Install]
WantedBy=multi-user.target
